# supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[inet_http_server]          ; inet (TCP) server disabled by default
port=127.0.0.1:9101         ; (ip_address:port specifier, *:port for all iface)

[supervisord]
logfile=%(ENV_LUNA_STUDIO_LOG_PATH)s/supervisord.log  ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid  ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
childlogdir=%(ENV_LUNA_STUDIO_LOG_PATH)s ; ('AUTO' child log dir, default $TEMP)
environment=LUNA_LIBS_PATH="%(ENV_LUNA_STUDIO_CONFIG_PATH)s/env"

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:ws-connector]
command=%(ENV_LUNA_STUDIO_BACKEND_PATH)s/ws-connector -v5
environment=LD_LIBRARY_PATH=%(ENV_OLD_LIBPATH)s
directory=%(ENV_PWD)s/
redirect_stderr=true

[program:broker]
command=%(ENV_LUNA_STUDIO_BACKEND_PATH)s/broker -v5
environment=LD_LIBRARY_PATH=%(ENV_OLD_LIBPATH)s
directory=%(ENV_PWD)s/
redirect_stderr=true

[program:luna-empire]
command=%(ENV_LUNA_STUDIO_BACKEND_PATH)s/luna-empire -v5
environment=LD_LIBRARY_PATH=%(ENV_OLD_LIBPATH)s
directory=%(ENV_PWD)s/
redirect_stderr=true
; stdout_logfile=%(ENV_LUNA_STUDIO_LOG_PATH)s/luna-empire.log
; stdout_logfile_maxbytes=1MB
; stdout_logfile_backups=10
; stdout_capture_maxbytes=1MB
; stdout_events_enabled=false
; stderr_logfile=luna-empire.err
; stderr_logfile_maxbytes=1MB
; stderr_logfile_backups=10
; stderr_capture_maxbytes=1MB

[program:logger]
command=%(ENV_LUNA_STUDIO_BACKEND_PATH)s/bus-logger -v4
environment=LD_LIBRARY_PATH=%(ENV_OLD_LIBPATH)s
directory=%(ENV_PWD)s/
redirect_stderr=true

[program:undo-redo]
command=%(ENV_LUNA_STUDIO_BACKEND_PATH)s/undo-redo -v5
environment=LD_LIBRARY_PATH=%(ENV_OLD_LIBPATH)s
directory=%(ENV_PWD)s/
redirect_stderr=true
; stdout_logfile=%(ENV_LUNA_STUDIO_LOG_PATH)s/undo-redo.log
; stdout_logfile_maxbytes=1MB
; stdout_logfile_backups=10
; stdout_capture_maxbytes=1MB
; stdout_events_enabled=false
; stderr_logfile=undo-redo.err
; stderr_logfile_maxbytes=1MB
; stderr_logfile_backups=10
; stderr_capture_maxbytes=1MB
