# supervisor/supervisord-package.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=%(ENV_LUNA_STUDIO_LOG_PATH)s/supervisord.log  ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid  ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
childlogdir=%(ENV_LUNA_STUDIO_LOG_PATH)s ; ('AUTO' child log dir, default $TEMP)
environment=LUNA_LIBS_PATH="%(ENV_LUNA_STUDIO_CONFIG_PATH)s/env",LD_LIBRARY_PATH="%(ENV_OLD_LIBPATH)s"

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:luna-atom]
environment=ATOM_HOME="%(ENV_LUNA_STUDIO_GUI_CONFIG_PATH)s"
command=%(ENV_LUNA_STUDIO_GUI_PATH)s -w %(ENV_LUNA_STUDIO_ATOM_ARG)s
priority=2
process_name=luna-atom
numprocs=1
startsecs=0
autorestart=false
stopasgroup=true
killasgroup=true
stdout_logfile=%(ENV_LUNA_STUDIO_LOG_PATH)s/luna-atom.log
stdout_logfile_maxbytes=0
stderr_logfile=%(ENV_LUNA_STUDIO_LOG_PATH)s/luna-atom.err
stderr_logfile_maxbytes=0

[eventlistener:luna-atom_exit]
command=%(ENV_LUNA_STUDIO_KILL_PATH)s
process_name=luna-atom-kill
events=PROCESS_STATE_EXITED

[program:luna-ws-connector]
environment=LD_LIBRARY_PATH=%(ENV_LUNA_STUDIO_BACKEND_LD_LIBRARY_PATH)s
command=%(ENV_LUNA_STUDIO_BACKEND_PATH)s/luna-ws-connector -v5
directory=%(ENV_PWD)s/
redirect_stderr=true

[program:luna-broker]
environment=LD_LIBRARY_PATH=%(ENV_LUNA_STUDIO_BACKEND_LD_LIBRARY_PATH)s
command=%(ENV_LUNA_STUDIO_BACKEND_PATH)s/luna-broker -v5
directory=%(ENV_PWD)s/
redirect_stderr=true

[program:luna-empire]
environment=LD_LIBRARY_PATH=%(ENV_LUNA_STUDIO_BACKEND_LD_LIBRARY_PATH)s
command=%(ENV_LUNA_STUDIO_BACKEND_PATH)s/luna-empire -v5
directory=%(ENV_PWD)s/
redirect_stderr=true
; stdout_logfile=%(ENV_LUNA_STUDIO_LOG_PATH)s/luna-empire.log
; stdout_logfile_maxbytes=1MB
; stdout_logfile_backups=10
; stdout_capture_maxbytes=1MB
; stdout_events_enabled=false
; stderr_logfile=luna-empire.err
; stderr_logfile_maxbytes=1MB
; stderr_logfile_backups=10
; stderr_capture_maxbytes=1MB

[program:luna-bus-logger]
environment=LD_LIBRARY_PATH=%(ENV_LUNA_STUDIO_BACKEND_LD_LIBRARY_PATH)s
command=%(ENV_LUNA_STUDIO_BACKEND_PATH)s/luna-bus-logger -v4
directory=%(ENV_PWD)s/
redirect_stderr=true

[program:luna-undo-redo]
environment=LD_LIBRARY_PATH=%(ENV_LUNA_STUDIO_BACKEND_LD_LIBRARY_PATH)s
command=%(ENV_LUNA_STUDIO_BACKEND_PATH)s/luna-undo-redo -v5
directory=%(ENV_PWD)s/
redirect_stderr=true
; stdout_logfile=%(ENV_LUNA_STUDIO_LOG_PATH)s/luna-undo-redo.log
; stdout_logfile_maxbytes=1MB
; stdout_logfile_backups=10
; stdout_capture_maxbytes=1MB
; stdout_events_enabled=false
; stderr_logfile=luna-undo-redo.err
; stderr_logfile_maxbytes=1MB
; stderr_logfile_backups=10
; stderr_capture_maxbytes=1MB
